# modules/walle/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:dind]
command=/bin/bash /usr/local/bin/dind-entrypoint
redirect_stderr=true

[program:sshd]
command=/usr/sbin/sshd -D
redirect_stderr=true

